# deb/mint-gnome-core/etc/xdg/autostart/mint-gnome-welcome.desktop
[Desktop Entry]
Name=Welcome
Icon=mintwelcome
Exec=mint-gnome-welcome
Terminal=false
Type=Application
