# packages/mint-gnome-welcome/etc/xdg/autostart/mint-gnome-welcome.desktop
[Desktop Entry]
Name=Welcome
Icon=lmg-os-logo-icon
Exec=mint-gnome-welcome
Terminal=false
Type=Application
OnlyShowIn=GNOME;
